# scripts/systemd/musu-bridge.service
[Unit]
Description=MUSU Bridge API Server
After=network.target

[Service]
Type=simple
# Placeholders __MUSU_ROOT__ are substituted by install-musu-bridge-service.sh
# when the file is copied into ~/.config/systemd/user/. Do NOT symlink — copy.
WorkingDirectory=__MUSU_ROOT__

# Load per-user config. Create it from scripts/systemd/bridge.env.example
EnvironmentFile=-%h/.musu/bridge.env

ExecStart=__MUSU_ROOT__/scripts/start-bridge.sh
Restart=on-failure
RestartSec=5

# Kill claude subprocesses aggressively on stop/restart.
# mixed = SIGTERM to main, SIGKILL to remaining after TimeoutStopSec.
# Without this, claude CLI ignores SIGTERM and becomes zombie.
KillMode=mixed
TimeoutStopSec=15
SendSIGKILL=yes

# Resource guardrails (tune via systemctl --user edit musu-bridge)
CPUQuota=50%
MemoryMax=2048M

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
