# scripts/systemd/musu-bee.service
[Unit]
Description=MUSU Bee Web UI (Next.js :3001)
After=network.target musu-bridge.service
Wants=musu-bridge.service

[Service]
Type=simple
WorkingDirectory=%h/musu-functions/musu-bee
ExecStartPre=/bin/bash -c 'test -d .next || npm run build || true'
ExecStart=/bin/bash -c 'if test -d .next; then exec node_modules/.bin/next start -p 3001; else exec node_modules/.bin/next dev -p 3001; fi'
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
