# scripts/systemd/musu-backup.timer
[Unit]
Description=Run MUSU SQLite backup daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
