# scripts/systemd/musud.service
[Unit]
Description=MUSU Process Supervisor (musud)
After=network.target

[Service]
Type=simple
# wiki/496 D12: was %h/musu-functions/... (Python-era layout). R6
# installer writes to ~/.musu/bin/ and the runtime resolves via dirs::home_dir().
WorkingDirectory=%h/.musu
ExecStart=%h/.musu/bin/musud
Environment=MUSU_HOME=%h/.musu
KillMode=mixed
TimeoutStopSec=20
SendSIGKILL=yes
Restart=on-failure
RestartSec=5

# Resource guardrails (same as musu-bridge)
CPUQuota=50%
MemoryMax=2048M

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
