# scripts/launchd/com.musu.bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.musu.bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__MUSU_ROOT__/scripts/start-bridge.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__MUSU_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Load env file (installer rewrites this on copy) -->
  <key>StandardOutPath</key>
  <string>__HOME__/.musu/logs/musu-bridge.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.musu/logs/musu-bridge.err.log</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <!-- Cap CPU and memory to avoid runaway claude subprocesses -->
  <key>SoftResourceLimits</key>
  <dict>
    <key>CPU</key>
    <integer>3600</integer>
  </dict>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
