# scripts/windows/musud_task.xml.tmpl
<?xml version="1.0" encoding="UTF-16"?>
<!--
  musu Scheduled Task template (wiki/496 §3 + D3 + F5/F19).

  Substitutes:
    {USER_ID}    — DOMAIN\USER of the operator
    {MUSU_HOME}  — operator's %USERPROFILE%\.musu directory
    {MUSUD_PATH} — full path to musud.exe

  Installer (install/platform/windows.rs) writes the rendered XML next
  to ~/.musu/install/musud_task.xml and runs:
    schtasks /Create /TN Musu\musud /XML <path> /F

  D3 hard constraints:
    - MultipleInstancesPolicy = StopExisting  (avoids silent boot failure
      when a zombie task lingers — IgnoreNew would refuse restart-on-logon)
    - RestartOnFailure Count=5 Interval=PT5M (less aggressive than the
      seed plan's PT1M/Count=3; D3 recalibration)
    - LogonType = InteractiveToken (no admin elevation; F5/F19)
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>MUSU process supervisor (musud) — logon-triggered.</Description>
    <URI>\Musu\musud</URI>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <UserId>{USER_ID}</UserId>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>{USER_ID}</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>StopExisting</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
    <RestartOnFailure>
      <Interval>PT5M</Interval>
      <Count>5</Count>
    </RestartOnFailure>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>{MUSUD_PATH}</Command>
      <WorkingDirectory>{MUSU_HOME}</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
